# wealthfolio/rootfs/etc/services.d/wealthfolio/finish
#!/usr/bin/execlineb -S0
# Take down the S6 supervision tree when wealthfolio-server stops.
/run/s6/basedir/bin/halt
