# aider/rootfs/etc/services.d/ttyd/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: Aider AI Coding Assistant
# Take down the S6 supervision tree when ttyd stops
# ==============================================================================
/run/s6/basedir/bin/halt
